# lib/systemd/system/lernstick-backup.service
[Unit]
Description=Background backup

[Service]
ExecStart=/usr/bin/lernstick_backup

[Install]
WantedBy=multi-user.target
